# systemd/netfilter-persistent.service.d/override.conf
[Service]
ExecStartPre=/etc/init.d/ipset start
ExecStopPost=/etc/init.d/ipset stop
